# scripts/schedule/com.redmountain.newhigh.wechat-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.redmountain.newhigh.wechat-daily</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>/Users/apple/Ahope/newhigh/news_collector_x_tweet.py</string>
		<key>--twitter-limit</key>
		<string>0</string>
		<key>--wechat-limit</key>
		<string>20</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/apple/Ahope/newhigh</string>
	<key>StartCalendarInterval</key>
	<dict>
		<!-- 每天检查一次微信公众号 (上午 9 点) -->
		<key>Hour</key><integer>9</integer>
		<key>Minute</key><integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/apple/Ahope/newhigh/logs/wechat_daily_stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/apple/Ahope/newhigh/logs/wechat_daily_stderr.log</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
