# integrations/hongshan/policy-news/com.newhigh.news-api.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.newhigh.news-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>REPLACE_WITH_NEWHIGH_ROOT/integrations/hongshan/policy-news/news_database.py</string>
        <string>api</string>
    </array>
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_NEWHIGH_ROOT/integrations/hongshan/policy-news</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>REPLACE_WITH_NEWHIGH_ROOT/integrations/hongshan/logs/news-api.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_NEWHIGH_ROOT/integrations/hongshan/logs/news-api-error.log</string>
</dict>
</plist>
